# httpd_container.service
# container-5f5ef8a03796fb6424492eb50849d83af502dced2077af31a6b0b82df534a375.service
# autogenerated by Podman 2.2.1
# Sat Apr  3 20:28:31 EEST 2021

[Unit]
Description=Podman container-5f5ef8a03796fb6424492eb50849d83af502dced2077af31a6b0b82df534a375.service
Documentation=man:podman-generate-systemd(1)
Wants=network.target
After=network-online.target

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=on-failure
ExecStart=/usr/bin/podman start 5f5ef8a03796fb6424492eb50849d83af502dced2077af31a6b0b82df534a375
ExecStop=/usr/bin/podman stop -t 10 5f5ef8a03796fb6424492eb50849d83af502dced2077af31a6b0b82df534a375
ExecStopPost=/usr/bin/podman stop -t 10 5f5ef8a03796fb6424492eb50849d83af502dced2077af31a6b0b82df534a375
PIDFile=/var/run/containers/storage/overlay-containers/5f5ef8a03796fb6424492eb50849d83af502dced2077af31a6b0b82df534a375/userdata/conmon.pid
KillMode=none
Type=forking

[Install]
WantedBy=multi-user.target default.target
